# scripts/cpu-tuning.service
[Unit]
Description=Unlock max boost for timekeeper VDF (governor=performance, EPP=performance, intel HWP dynamic boost)
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# `-` prefix tolerates per-ExecStart failures (e.g. paths that don't exist on AMD or kernels without intel_pstate).
# The inner `[ -w PATH ]` check avoids writing where the path is missing so we don't spam errors in the journal.
ExecStart=-/bin/bash -c 'for f in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do [ -w "$f" ] && echo performance > "$f"; done'
ExecStart=-/bin/bash -c 'for f in /sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference; do [ -w "$f" ] && echo performance > "$f"; done'
ExecStart=-/bin/bash -c '[ -w /sys/devices/system/cpu/intel_pstate/hwp_dynamic_boost ] && echo 1 > /sys/devices/system/cpu/intel_pstate/hwp_dynamic_boost || true'

[Install]
WantedBy=multi-user.target
